# systemd/check-battery.service
[Unit]
Description=Vérification du niveau de batterie

[Service]
ExecStart=/home/$USER/.local/bin/check-battery.sh
